# plugins/claude-code-hermit/state-templates/watchdog/hermit-watchdog@.service
[Unit]
Description=Hermit watchdog — dead-session restart + wedge detection for {{NAME}}
After=network.target

[Service]
Type=oneshot
WorkingDirectory={{ROOT}}
ExecStart={{ROOT}}/.claude-code-hermit/bin/hermit-watchdog run
StandardOutput=journal
StandardError=journal
